# systemd/array-firewall-fleet-sync.service
[Unit]
Description=Array Firewall fleet blocklist sync
After=array-firewall-api.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/array-firewall/scripts/fleet-sync.py
Nice=10
