# systemd/array-firewall-ai-ops.timer
[Unit]
Description=AI Autopilot fusion cycle every 30 seconds

[Timer]
OnBootSec=45s
OnUnitActiveSec=30s
Persistent=true

[Install]
WantedBy=timers.target
