# systemd/array-firewall-ai-ops.timer
[Unit]
Description=AI Autopilot fusion cycle every 5 seconds

[Timer]
OnBootSec=15s
OnUnitActiveSec=5s
Persistent=true

[Install]
WantedBy=timers.target
